# task5/ansible/roles/maintenance/templates/docker-prune.service.j2
[Unit]
Description=Docker Cleanup

[Service]
Type=oneshot
ExecStart=/usr/bin/docker system prune -af --volumes
