# task5/ansible/roles/maintenance/templates/docker-prune.timer.j2
[Unit]
Description=Run Docker cleanup daily

[Timer]
OnCalendar=*-*-* {{ docker_prune_time }}
Persistent=true

[Install]
WantedBy=timers.target
