# task4/ansible/roles/maintenance/templates/prometheus-backup.timer.j2
[Unit]
Description=Backup Prometheus data on schedule

[Timer]
OnCalendar=*-*-* {{ prometheus_backup_times }}
Persistent=true

[Install]
WantedBy=timers.target
